# UniCron/UniCron.cron
*/5 * * * * /data/adb/modules/UniCron/UniCron.sh